# theme-project/news-service/parser/crontab
0 8 * * * python /app/parser/news_parser.py >> /var/log/parser.log 2>&1
